# data/cd-annual-converter-master/vm/mailhog.service
[Unit]
Description=Mailhog
After=network.service

[Service]
ExecStart=/opt/mailhog > /dev/null 2>&1 &

[Install]
WantedBy=multi-user.target
